# systemd/downloadupdate.service
[Unit]
Description=Download system updates

[Service]
Type=simple
ExecStart=pacman -Syyuw
